# systemd/nimo-ping.service
[Unit]
Description=running speed test

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nimo ping
