# systemd/nimo-speed-test.service
[Unit]
Description=running ping test

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nimo speed-test
